# .misc/macos/launchd/Library.home/LaunchAgents/homebrew.mxcl.mutagen.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>homebrew.mxcl.mutagen</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>MUTAGEN_DATA_DIRECTORY</key>
		<string>$MUTAGEN_DATA_DIRECTORY</string>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/mutagen</string>
		<string>daemon</string>
		<string>run</string>
	</array>
	<key>StandardOutPath</key>
	<string>/opt/homebrew/var/log/mutagen.log</string>
	<key>StandardErrorPath</key>
	<string>/opt/homebrew/var/log/mutagen.log</string>
	<key>RunAtLoad</key>
	<true/>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
		<string>Background</string>
		<string>LoginWindow</string>
		<string>StandardIO</string>
		<string>System</string>
	</array>
</dict>
</plist>
